# conf/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   
chmod=0700                       

[supervisord]
logfile=/var/log/supervisor/supervisord.log 
pidfile=/var/run/supervisord.pid 
childlogdir=/var/log/supervisor
nodaemon=true  

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 

[program:elasticsearch]
command=JAVA_OPTIONS=-Xmx512M /usr/share/elasticsearch/bin/elasticsearch
autorestart=true

[program:kibana]
command=/opt/kibana/kibana-4.6.1-linux-x86_64/bin/kibana -c /etc/kibana/kibana.yml
autorestart=true

[program:logstash]
command=/opt/logstash/bin/logstash -f /etc/logstash/conf.d/logstash.conf
autorestart=true
